# dot_config/systemd/user/t3code.service
[Unit]
Description=T3 Code headless web UI
Documentation=https://github.com/pingdotgg/t3code

[Service]
Type=simple
WorkingDirectory=%h
ExecStart=/usr/bin/npm exec --yes --package t3@0.0.28 -- t3 serve --mode web --host 127.0.0.1 --port 3773 --base-dir %h/.local/share/t3code --no-browser --tailscale-serve --tailscale-serve-port 443 %h
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
