# dot_config/systemd/user/niri.service.d/override.conf.tmpl
{{ if eq .chezmoi.hostname "aeolus" -}}
[Service]
ExecStart=
ExecStart=%h/.local/share/cargo/bin/niri --session
{{ end -}}
